# kk.conf
[program:kk]
command=python3 /home/curator2/kk/main.py
autostart=true
autorestart=true
stderr_logfile=/home/curator2/kk/test.err.log
stdout_logfile=/home/curator2/kk/test.out.log
